# etc/colossus_gunicorn_supervisord.conf
[program:colossus_gunicorn]
command=/home/colossus/colossus/bin/gunicorn_start
user=colossus
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/colossus/logs/gunicorn.log
